# resources/systemd/btrfs-backup.service
[Unit]
Description=Take BTRFS snapshots
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/BtrfsBackup
